# sw/vm/overlay/etc/profile.d/env.sh
HOSTNAME="$(cat /etc/hostname)"
export HOSTNAME

HOME=/root
export HOME

. /root/.bashrc
